# .config/systemd/user/wallpapers-remount.service
[Unit]
Description=Restart wallpaper mount if needed

[Service]
ExecStart=/home/aujvary/src/sh/restart-mount.sh

